# run_at_launch.conf
description "run_at_launch.conf"
author      "JR King"

start on startup

stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    HOST=`hostname`
    export HOME="/home/ubuntu"
    chdir /home/ubuntu/
    exec /bin/sh /home/ubuntu/aws/update.sh
    echo "update started"
end script
